# openstatus-client/DEBIAN/postinst
#!/bin/sh
set -e
update-rc.d openstatus-client defaults
service openstatus-client restart
